# root/etc/uci-defaults/luci-haproxy-tcp
#!/bin/sh

/etc/init.d/haproxy-tcp enable

exit 0
